# systemd/cpuls.service
[Unit]
Description=CPU Log Service

[Service]
ExecStart=/home/dosx/repos/cli/cpuls/target/debug/cpuls

[Install]
WantedBy=multi-user.target
